# files/reddy.service
[Unit]
Description=Reddy

[Service]
ExecStart=/etc/reddy/reddy
Restart=always