# files/reddy.service
[Unit]
Description=Reddy

[Service]
WorkingDirectory=/etc/reddy
ExecStart=nohup ./reddy &
Restart=always